# launchd/com.nightclaude.pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- macOS equivalent of nightclaude-pull.service: fetch overnight
         results from the worker at login. __HOME__ is filled in by
         install.sh (launchd does not expand ~ or $HOME). -->
    <key>Label</key>
    <string>com.nightclaude.pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <!-- Give the network / the worker's mDNS a moment after login. -->
        <string>sleep 15; exec "__HOME__/.local/bin/nightclaude" pull</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.local/share/nightclaude/logs/launchd-pull.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/share/nightclaude/logs/launchd-pull.log</string>
</dict>
</plist>
